# version2/mdns-to-dns/mdns-to-dns.service
[Unit]
Description=CO2 sensor mDNS to DNS
After=network.target

[Service]
Type=oneshot
ExecStart=/home/user/.local/bin/mdns-to-dns.sh

[Install]
WantedBy=default.target
