# version2/mdns-to-dns/mdns-to-dns.timer
[Unit]
Description=Run mdns-to-dns.sh every 5 minutes

[Timer]
OnBootSec=0min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
